# tailscale/rootfs/etc/s6-overlay/s6-rc.d/share-services/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community App: Tailscale
# Advertises configured Tailscale Services from this host
# ==============================================================================

shopt -s extglob

declare current_config
declare service
declare service_name
declare protocol
declare port
declare target
declare path

declare -a options=()
declare -a configured_services=()
declare -a configured_names=()
declare -a current_names=()

declare https_checked=false

# Read configured services
# Note: bashio::config emits one JSON object per line for array options
readarray -t configured_services < <(bashio::config 'services')

# Build a list of configured service names
for service in "${configured_services[@]}"; do
    name=$(bashio::jq "${service}" '.name')
    configured_names+=("${name}")
done

if bashio::var.has_value "${configured_names[*]}"; then
    bashio::log.info "Configured service hosts: ${configured_names[*]}"
fi

# Read currently advertised services from the tailscale serve config
if ! current_config=$(/opt/tailscale serve get-config --all 2>/dev/null); then
    current_config='{}'
fi

# Build a list of currently configured service names
readarray -t current_names < <(bashio::jq "${current_config}" '.services // {} | keys[]' || true)

# Remove services that are no longer configured in this app
for service_name in "${current_names[@]}"; do
    if ! printf '%s\n' "${configured_names[@]}" | grep -qx "${service_name}"; then
        bashio::log.info "Removing stale service host configuration for ${service_name}"
        if ! /opt/tailscale serve drain "${service_name}" >/dev/null 2>&1; then
            bashio::log.warning "Failed to drain stale service host ${service_name}"
        fi
        if ! /opt/tailscale serve clear "${service_name}" >/dev/null 2>&1; then
            bashio::log.warning "Failed to clear stale service host ${service_name}"
        fi
    fi
done

# Configure and advertise each service
for service in "${configured_services[@]}"; do
    service_name=$(bashio::jq "${service}" '.name')
    protocol=$(bashio::jq "${service}" '.protocol')
    port=$(bashio::jq "${service}" '.port')
    target=$(bashio::jq "${service}" '.target')
    path=$(bashio::jq "${service}" '.path // empty')

    if ! [[ "${protocol}" = @(http|https|tcp|tls-terminated-tcp) ]]; then
        bashio::exit.nok "Invalid value '${protocol}' for services.protocol"
    fi

    if ! "${https_checked}" && [[ "${protocol}" = @(https|tls-terminated-tcp) ]]; then
        if ! tailscale_status=$(/opt/tailscale status --self=true --peers=false --json 2>/dev/null); then
            bashio::exit.nok "Unable to determine Tailscale HTTPS support"
        fi
        if [[ "$(bashio::jq "${tailscale_status}" '.Self.CapMap | has("https")')" != "true" ]]; then
            bashio::exit.nok "Tailscale's HTTPS support is disabled"
        fi
        https_checked=true
    fi

    options=("serve")
    options+=("--service=${service_name}")
    options+=("--${protocol}=${port}")

    if bashio::var.has_value "${path}"; then
        options+=("--set-path=${path}")
    fi

    options+=("${target}")

    # Clean any existing configuration for this service before reconfiguring
    if ! /opt/tailscale serve drain "${service_name}" >/dev/null 2>&1; then
        bashio::log.warning "Failed to drain service host ${service_name}"
    fi
    if ! /opt/tailscale serve clear "${service_name}" >/dev/null 2>&1; then
        bashio::log.warning "Failed to clear service host ${service_name}"
    fi

    bashio::log.info "Advertising service host for ${service_name}"
    if ! /opt/tailscale "${options[@]}"; then
        bashio::exit.nok "Failed to advertise service host for ${service_name}"
    fi
done

exit 0
